# upstart_CltCodeCamp2014Notify.conf
description "Notification script to run at shutdown"

start on starting rc
task
exec /usr/bin/curl http://localhost:8080/api/notifyShutdown